# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
Documentation=man:nix-daemon https://bsdos.org/manual
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
RequiresMountsFor=@localstatedir@/bsd/db
ConditionPathIsReadWrite=@localstatedir@/bsd/daemon-socket

[Service]
ExecStart=@@bindir@/bsd-daemon nix-daemon --daemon
KillMode=process
LimitNOFILE=1048576
TasksMax=1048576
Delegate=yes

[Install]
WantedBy=multi-user.target
