# misc/systemd/nix-daemon.socket.in
[Unit]
Description=Nix Daemon Socket
Before=multi-user.target
RequiresMountsFor=@storedir@
ConditionPathIsReadWrite=@localstatedir@/bsd/daemon-socket

[Socket]
ListenStream=@localstatedir@/bsd/daemon-socket/socket

[Install]
WantedBy=sockets.target
